# roles/docker-service/templates/docker.service.j2
#***************************************************
#Description: this file is generated by ansible
#             and system should use this to detect some state
#             so do not modify it !!!
#info:{{ imageName }}:{{ imageTag }}
#***************************************************

[Unit]
Description=Daemon for {{ serviceName }}
After=docker.service
Wants=
Requires=docker.service

[Service]
Restart=on-failure
StartLimitInterval=20
StartLimitBurst=5
TimeoutStartSec=0
Environment="HOME=/root"
ExecStartPre=-/usr/bin/docker kill {{ serviceName }}
ExecStartPre=-/usr/bin/docker rm  {{ serviceName}}

ExecStart=/usr/bin/docker run --name {{ serviceName }} \
          {% for item in tcp %}
          -p {{item}}/tcp \
          {% endfor %}
          {% for item in udp %}
          -p {{item}}/udp \
          {% endfor %}
          {% for item in volumns %}
          -v {{item}} \
          {% endfor %}
          {{ imageName }}:{{ imageTag }} {{command}}
         
ExecStop=-/usr/bin/docker stop {{ serviceName }}
ExecStop=-/usr/bin/docker rm  {{ serviceName }}

[Install]
WantedBy=multi-user.target

